# content/etc/s6-overlay/s6-rc.d/init-adduser/run
#!/command/with-contenv bash

echo "Checking for abc user"
PUID=${PUID:-3001}
PGID=${PGID:-3000}

# Keep startup non-fatal if user/group already exists or cannot be recreated.
/usr/sbin/addgroup -g "${PGID}" abc 2>/dev/null || echo "abc group already exists or could not be created"
/usr/sbin/adduser -D -u "${PUID}" -G abc -h /app -g "App User" abc 2>/dev/null || echo "abc user already exists or could not be created"

echo "
───────────────────────────────────────
User UID:    $(id -u abc)
User GID:    $(id -g abc)
───────────────────────────────────────
"

for directory in /app /config /defaults; do
  if [ ! -d "${directory}" ]; then
    echo "Creating directory: ${directory}"
    mkdir -p "${directory}"
  fi

  echo "Ensuring ownership for: ${directory}"
  chown -R abc:abc "${directory}"
done
