# content/etc/s6-overlay/s6-rc.d/init-tz/run
#!/command/with-contenv bash

TZ=${TZ:-America/New_York}
ZONEINFO_PATH="/usr/share/zoneinfo/${TZ}"

echo "Setting timezone to $TZ"

if [ ! -e "${ZONEINFO_PATH}" ]; then
  echo "Invalid timezone '${TZ}', falling back to America/New_York"
  TZ="America/New_York"
  ZONEINFO_PATH="/usr/share/zoneinfo/${TZ}"
fi

if [ ! -f /etc/timezone ] || ! grep -qxF "$TZ" /etc/timezone; then
  echo "$TZ" > /etc/timezone
  ln -sf "${ZONEINFO_PATH}" /etc/localtime
fi